# etc/systemd/user/shrimpos.service
[Unit]
Description=ShrimpOS-Service

[Service]
Type=simple
ExecStart=/etc/shrimpos/shrimpOS-startup-system.sh

[Install]
WantedBy=default.target
