# etc/systemd/system/shrimpos.service
[Unit]
Description=ShrimpOS-Service

[Service]
Type=simple
ExecStart=/etc/shrimpos/system-service.sh

[Install]
WantedBy=default.target
